# dnscontrol.crontab
30 7 * * *  install_path/dnscontrol/run.sh update >> install_path/dnscontrol/dnscontrol.log 2>&1
